# deploy/inapp-bot.service
[Unit]
Description=inApp Telegram bot (@inAppProBot)
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/opt/badcomment
Environment=NODE_ENV=production
Environment=DATABASE_URL=file:/opt/badcomment/data/prod.db
# TELEGRAM_BOT_TOKEN lives only in the gitignored prod env (chmod 600), never in git.
EnvironmentFile=/opt/badcomment/.env
ExecStart=/usr/bin/node bot/bot.mjs
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
